# wg-quick-auto-restart@.path
[Unit]
Description=Watch a WireGuard config fragment directory for changes and restarts wg-quick

[Path]
DirectoryNotEmpty=/etc/wireguard/%i
PathChanged=/etc/wireguard/%i
Unit=wg-quick-auto-restart@%i.service

[Install]
WantedBy=paths.target

